# production/modules/pt_setup/templates/psft-db.erb
#!/bin/bash
#
# chkconfig: 35 94 16
# description: This is a program that is responsible for taking care of
# configuring the Oracle Database 12g and its associated
# services.
#

if [ -f /etc/init.d/functions ]; then
    . /etc/init.d/functions
elif [ -f /etc/rc.d/init.d/functions ]; then
    . /etc/rc.d/init.d/functions
else
    exit 1
fi
unset TMPDIR

SCRIPT_NAME=$(basename $(readlink -nf $0))
LOG_FILE=/var/log/$SCRIPT_NAME.log
CONTAINER_NAME=<%= @container_name %>
DATABASE_NAME=<%= @db_name %>
ORACLE_OWNER=<%= @oracle_user %>
ORACLE_HOME=<%= @oracle_home_dir %>
ORACLE_SID=$CONTAINER_NAME
LD_LIBRARY_PATH=$ORACLE_HOME/lib:$LD_LIBRARY_PATH
PATH=$ORACLE_HOME/bin:$PATH

THIS_OS="THIS_OS-not-initialized"
THIS_OS=`uname -s`

log_message() {
    echo "[LOG]  $(date +'%b %d %H:%M:%S') ${0##*/}: $*" >> $LOG_FILE
}

get_listener_name() {
  listener_name=$(cat $ORACLE_HOME/network/admin/listener.ora | grep "^[a-zA-Z].*=$" | cut -f1 -d '=')
  echo $listener_name
}

# function start listener
start_listener () {
    export ORACLE_HOME
    prog="PeopleSoft Database Listener"
    echo -n $"Starting $prog: "

    LISTENER_NAME=$(get_listener_name)
    log_message "Starting $prog"

    case "$THIS_OS" in
       SunOS ) su - $ORACLE_OWNER -c "$ORACLE_HOME/bin/lsnrctl start $LISTENER_NAME"            1>>$LOG_FILE 2>&1
               ;;
           * ) su -s /bin/bash $ORACLE_OWNER -c "$ORACLE_HOME/bin/lsnrctl start $LISTENER_NAME" 1>>$LOG_FILE 2>&1
               ;;
    esac
    ret=$?
    [ $ret -eq 0 ] && success || failure
    echo
    log_message "Started $prog"

    return $ret
}

# function stop listener
stop_listener () {
    ret=0
    export ORACLE_HOME
    prog="PeopleSoft Database Listener"
    echo -n $"Stopping $prog: "

    LISTENER_NAME=$(get_listener_name)
    log_message "Stopping $prog"
    # check if the listener is already stopped
    status=`ps -ef | grep tns | grep $LISTENER_NAME`
    if [ "$status" == "" ]; then
        log_warn "$prog has already been stopped" >/dev/null 2>&1
        success
        ret=1
    else
        case "$THIS_OS" in
           SunOS ) su - $ORACLE_OWNER -c "$ORACLE_HOME/bin/lsnrctl stop $LISTENER_NAME"            1>>$LOG_FILE 2>&1
                   ;;
               * ) su -s /bin/bash $ORACLE_OWNER -c "$ORACLE_HOME/bin/lsnrctl stop $LISTENER_NAME" 1>>$LOG_FILE 2>&1
                   ;;
        esac
        ret=$?
        [ $ret -eq 0 ] && success || failure
    fi
    echo
    log_message "Stopped $prog"

    return $ret
}

start_db_container() {
    export ORACLE_SID ORACLE_HOME
    export TNS_ADMIN=$ORACLE_HOME/network/admin

    prog="PeopleSoft Container Database $CONTAINER_NAME"
    echo -n $"Starting $prog: "

    log_message "Starting $prog: "
    su -s /bin/bash  $ORACLE_OWNER -c "$ORACLE_HOME/bin/sqlplus -s /nolog" <<SQL 2>&1>> $LOG_FILE
WHENEVER SQLERROR EXIT 1
connect / as sysdba
startup
SQL
    ret=$?
    [ $ret -eq 0 ] && success || failure
    echo
    log_message "Started $prog"

    return $ret
}

stop_db_container() {
    export ORACLE_SID ORACLE_HOME
    prog="PeopleSoft Container Database $CONTAINER_NAME"
    echo -n $"Stopping $prog: "

    log_message "Stopping $prog"
    su -s /bin/bash  $ORACLE_OWNER -c "$ORACLE_HOME/bin/sqlplus -s /nolog" <<SQL 2>&1>> $LOG_FILE
WHENEVER SQLERROR EXIT 1
connect / as sysdba
shutdown immediate
SQL
    ret=$?
    [ $ret -eq 0 ] && success || failure
    echo
    log_message "Stopped $prog"

    return $ret
}

get_database_status() {
    export ORACLE_SID ORACLE_HOME

    CHECK_PDB_SCRIPT=/tmp/checkpdb.sql

cat >> $CHECK_PDB_SCRIPT <<-EOF
set termout on
set echo off
set verify off
set feedback off
set heading off
select 'CDB' from v\$database where name='&1' and cdb='YES';
select 'PDBSEED' From v\$pdbs where name='PDB\$SEED' and open_mode = 'READ ONLY';
select '&2' , open_mode From v\$pdbs where name='&2' and open_mode in ('READ WRITE','MOUNTED');
exit
EOF
chmod 755 $CHECK_PDB_SCRIPT

    log_message "Getting the status of PDB $DATABASE_NAME"

    STATUS_FILE=/tmp/$DATABASE_NAME.status
    su -s /bin/bash  $ORACLE_OWNER -c "$ORACLE_HOME/bin/sqlplus -s /nolog" <<SQL 2>&1> $STATUS_FILE
connect / as sysdba
@$CHECK_PDB_SCRIPT $CONTAINER_NAME $DATABASE_NAME
SQL
rm -rf $CHECK_PDB_SCRIPT
}

start() {
    # start NET Listner last
    start_listener
    list_ret=$?

    # start the container
    start_db_container
    db_ret=$?

    final_ret=$(( list_ret | db_ret ))
    [ $final_ret -eq 0 ] && touch $LOCKFILE
    return $final_ret
}

stop() {

    # stop NET Listner first
    stop_listener
    list_ret=$?

    # stop the database
    stop_db_container
    db_ret=$?

    final_ret=$(( list_ret | db_ret ))
    [ $final_ret -eq 0 ] && rm -f $LOCKFILE
    return $final_ret
}

dostatus() {
    export ORACLE_SID ORACLE_HOME

    # check if the container database is running
    status=`ps -ef | grep ora_pmon_$CONTAINER_NAME | grep -v grep`
    if [ "$status" == "" ]; then
        echo "PeopleSoft Container Database $CONTAINER_NAME Status is Down"
    else
        echo "PeopleSoft Container Database $CONTAINER_NAME Status is Up"
        get_database_status
        DB_STATUS_FILE=/tmp/$DATABASE_NAME.status
        if [ -f $DB_STATUS_FILE ]; then
          if grep -q CDB $DB_STATUS_FILE && grep -q PDBSEED $DB_STATUS_FILE; then
            if grep -q $DATABASE_NAME $DB_STATUS_FILE; then
              if grep -q MOUNTED $DB_STATUS_FILE; then
                echo "PeopleSoft Pluggable Database $DATABASE_NAME Status is Closed"
              else
                echo "PeopleSoft Pluggable Database $DATABASE_NAME Status is Open"
              fi
            else
              echo "PeopleSoft Pluggable Database $DATABASE_NAME Does not Exists"
            fi
          else
            echo "PeopleSoft Pluggable Database $DATABASE_NAME Status is UNKNOWN"
          fi
          rm -rf $DB_STATUS_FILE >/dev/null
        else
          echo "PeopleSoft Pluggable Database $DATABASE_NAME Status is UNKNOWN"
        fi
    fi
    # check if the listener is running
    LISTENER_NAME=$(get_listener_name)
    status=`ps -ef | grep tns | grep $LISTENER_NAME`
    if [ "$status" == "" ]; then
        echo "PeopleSoft Database Listener is Down"
    else
        echo "PeopleSoft Database Listener is Up"
    fi
}

##### main ##############################################################################
# Set path if path not set (if called from /etc/rc)
case $PATH in
    "") PATH=/bin:/usr/bin:/sbin:/etc
        export PATH
        ;;
esac

export LC_ALL=C

if [ "$(id -u)" != "0" ]; then
    echo "You must be root to run this script."
    exit 1
fi

LOCKFILE=<%= @services_lock_dir %>/psft-db

# See how we were called
case "$1" in
    start)
      start
      ;;
    stop)
      stop
      ;;
    restart|reload)
      restart
        ;;
    status)
      dostatus
      ;;
    *)
      echo "Usage: $0 {start|stop|restart|status}"
      exit 1
esac
